# php72/lang/php72m/pkg-plist.fpm
etc/php/7.2/fpm/php-fpm.conf.default
etc/php/7.2/fpm/pool.d/www.conf.default
etc/rc.d/php7.2-fpm
man/man8/php-fpm7.2.8.gz
sbin/php-fpm7.2
sbin/php_fpm72
share/php/7.2/fpm/status.html
